# builds/install/arch-specific/linux/scratchbird.service.in
[Unit]
Description=ScratchBird Database Server
After=syslog.target network.target
Documentation=https://github.com/dcalford/ScratchBird

[Service]
User=scratchbird
Group=scratchbird
Type=forking
ExecStart=@FB_SBINDIR@/sb_fbguard -daemon -forever
Environment=SCRATCHBIRD=@FB_PREFIX@

[Install]
WantedBy=multi-user.target
